# nagaiku-budget-backend.service
[Unit]
Description=NPO Budget Management System - Backend (統一環境)
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=tanaka
Group=tanaka
WorkingDirectory=/home/tanaka/nagaiku-budget/backend
# 統一環境用環境変数設定
Environment=PORT=8000
Environment=DATABASE_NAME=nagaiku_budget
Environment=DATABASE_USER=nagaiku_user
Environment=DATABASE_PASSWORD=nagaiku_password2024
Environment=DATABASE_HOST=localhost
Environment=DATABASE_PORT=5432
ExecStart=/home/tanaka/nagaiku-budget/venv/bin/uvicorn main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=10
StandardOutput=append:/home/tanaka/nagaiku-budget/logs/backend_prod.log
StandardError=append:/home/tanaka/nagaiku-budget/logs/backend_prod.log

# プロセス管理
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target 